# rootfs/etc/s6-overlay/s6-rc.d/vdlm2dec/run
#!/bin/sh
exec /etc/s6-overlay/scripts/vdlm2dec
